# .config/awesome/autostart/flameshot.desktop
[Desktop Entry]
Name=Flameshot
Exec=/usr/bin/flameshot
NoDisplay=true
Type=Application
