# autostart/emacsdaemon.desktop
[Desktop Entry]
Name=Emacs Daemon
Exec=emacs --daemon
NoDisplay=true
Type=Application
