# autostart/eww.desktop
[Desktop Entry]
Name=eww
Exec=eww open-many bar status
NoDisplay=true
OnlyShowIn=Hyprland
Type=Application
