# autostart/eww.desktop
[Desktop Entry]
Name=eww
Exec=eww open bar
NoDisplay=true
OnlyShowIn=Hyprland
Type=Application
