# autostart/gammastep.desktop
[Desktop Entry]
Name=Gammastep
Exec=gammastep-indicator
NoDisplay=true
Type=Application
