# autostart/numlockx.desktop
[Desktop Entry]
Name=numlockx
Exec=/usr/bin/numlockx on
NoDisplay=true
Type=Application
