# autostart/setxkbmap.desktop
[Desktop Entry]
Name=setxkbmap
Exec=/usr/bin/setxkbmap -layout latam
NoDisplay=true
Type=Application
