# autostart/waybar.desktop
[Desktop Entry]
Name=Waybar
Exec=waybar
NoDisplay=true
OnlyShowIn=Hyprland
Type=Application
