# .config/hypr/autostart/eww.desktop
[Desktop Entry]
Name=eww
Exec=/usr/bin/eww open bar
NoDisplay=true
Type=Application
