# .config/hypr/autostart/wpaperd.desktop
[Desktop Entry]
Name=wpaperd
Exec=/usr/bin/wpaperd
NoDisplay=true
Type=Application
