# .config/qtile/autostart/dunst.desktop
[Desktop Entry]
Name=Dunst
Exec=$HOME/.dotfiles/scripts/reload_dunst.sh
NoDisplay=true
Type=Application
