# .config/qtile/autostart/dunst.desktop
[Desktop Entry]
Name=Dunst
Exec=/usr/bin/dunst
NoDisplay=true
Type=Application
