# .config/qtile/autostart/dunst.desktop
[Desktop Entry]
Name=Dunst
Exec=/bin/bash -c "~/.dotfiles/scripts/reload_dunst.sh"
NoDisplay=true
Type=Application
