# .config/qtile/autostart/picom.desktop
[Desktop Entry]
Name=Picom
Exec=env DISPLAY=":0" /usr/bin/picom --unredir-if-possible
NoDisplay=true
Type=Application
