# .config/qtile/autostart/picom.desktop
[Desktop Entry]
Name=Picom
Exec=env DISPLAY=":0" /usr/bin/picom
NoDisplay=true
Type=Application
